# streamScript/deploy/supervisor.conf
[program:streamScript]
command: /usr/bin/python -m wsgi
directory: /home/ubuntu/streamScript
autostart: true
stdout_logfile: /tmp/out_log
stderr_redirect: false
stderr_logfile: /tmp/err_log