# tweetTrack/deploy/supervisord.conf
[program:tweetTrack]
command: /usr/bin/python -m wsgi
directory: /home/ubuntu/tweetTrack/
autostart: true
stderr_redirect: true
stdout_logfile: /tmp/tweetTracklog