# tweetTrack/deploy/supervisord.conf
[program:twitter_stream]
command: /usr/bin/python twitter_stream.py
directory: /home/ubuntu/tweetTrack
autostart: true
stderr_redirect: true
stdout_logfile: /tmp/tweetTracklog